# ansible/monitoring/templates/node_exporter.service.j2
[Unit]
Description=simple node exporter service
After=network.target

[Service]
Type=simple
ExecStart=/opt/node_exporter/node_exporter --web.listen-address=:9110
ExecStop=killall node_exporter

[Install]
WantedBy=multi-user.target
